# home/kiosk/.config/autostart/kiosk.desktop
[Desktop Entry]
Type=Application
Exec=/home/kiosk/start-kiosk.sh
